# src/NHSS.QSDK_MIPS.4.0/apss_proc/out/proprietary/QSDK-Base/qsdk-luci-app-0.2.2922.227/qca/src/luci-app-atf/root/etc/uci-defaults/luci-app-atf
#!/bin/sh

# Copyright (c) 2016 Qualcomm Atheros, Inc.
#
# All Rights Reserved.
# Qualcomm Atheros Confidential and Proprietary
uci batch <<-EOF
	set luci.main.mediaurlbase=/luci-static/atf
	set luci.sauth.sessiontime=$((24*60*60))
	commit luci
EOF
exit 0
