# sw/qsdk/package/base-files/files/etc/init.d/set_upgrade_config
#!/bin/sh /etc/rc.common

START=99
start() {
       fw_setenv fw_type "production"
       fw_setenv customer_fw "no"
}